# ratinatorly.conf
description 'ratinatorly'
start on startup
chdir /home/node/ratinatorly
respawn
exec sudo -u node sh -c "node app.js >> log.txt 2>&1"
