# ratinatorly.conf
description 'ratinatorly'
start on startup
chdir /root/ratinatorly
respawn
exec sudo -u root sh -c "PORT=80 node app.js"
